# rootfs/etc/s6-overlay/s6-rc.d/svc-icecast/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid icecast2 icecast2 -c /etc/icecast2/icecast.xml
