# rootfs/etc/s6-overlay/s6-rc.d/svc-darkice/run
#!/usr/bin/with-contenv bash

export HOME="$( getent passwd ubuntu | cut -d: -f6 )"
exec s6-setuidgid ubuntu darkice -c /etc/darkice.cfg
